# formulas/dumping_formula/res/vjoule_dumping_formula.service
[Unit]
Description=vjoule dumping formula
Wants=vjoule_sensor.service

[Service]
User=root
WorkingDirectory=/etc/vjoule/
ExecStart=/usr/bin/vjoule_dumping_formula
Restart=always
Slice=vjoule.slice


[Install]
WantedBy=multi-user.target
